# plinth/modules/users/data/usr/lib/systemd/system/nslcd.service.d/freedombox.conf
[Service]
# Invalidate LDAP related caches.
ExecStartPost=-nscd --invalidate=passwd
ExecStartPost=-nscd --invalidate=group
