# plinth/modules/transmission/data/usr/lib/systemd/system/transmission-daemon.service.d/freedombox.conf
# The service depends on users and groups defined in LDAP
[Unit]
After=nslcd.service
