# plinth/modules/bepasty/data/usr/lib/systemd/system/uwsgi-app@bepasty-freedombox.service.d/freedombox.conf
[Service]
User=bepasty
Group=bepasty
StateDirectory=bepasty
LoadCredential=bepasty-freedombox.conf:/etc/bepasty-freedombox.conf
