# plinth/modules/minetest/data/usr/lib/systemd/system/minetest-server.service.d/freedombox.conf
# Workaround for Debian bug #1090067
[Service]
ExecStart=
ExecStart=/usr/lib/minetest/minetestserver --config /etc/minetest/minetest.conf --logfile /var/log/minetest/minetest.log --gameid minetest_game
