# plinth/modules/ttrss/data/usr/lib/systemd/system/tt-rss.service.d/freedombox.conf
# Restart the service every 120 seconds always. When tt-rss can't connect to a
# database temporarily, it will exist with exit code 101. 120 seconds is the
# default daemon sleep interval for tt-rss.
[Service]
Restart=always
RestartSec=120s
