# plinth/modules/radicale/data/usr/lib/systemd/system/uwsgi-app@radicale.service.d/freedombox.conf
[Service]
User=radicale
Group=radicale
StateDirectory=radicale
