# plinth/modules/miniflux/data/usr/lib/systemd/system/miniflux.service.d/freedombox.conf
# FreedomBox configuration file stores both static settings and user
# preferences. These settings are loaded as environment variables. Hence, they
# take precedence, overriding the settings in miniflux.conf.

[Service]
EnvironmentFile=/etc/miniflux/freedombox.conf
