# plinth/modules/shadowsocks/data/lib/systemd/system/shadowsocks-libev-local@.service.d/freedombox.conf
[Service]
StateDirectory=shadowsocks-libev/%i
DynamicUser=yes
