# plinth/modules/config/data/lib/systemd/system/zramswap.service.d/freedombox.conf
# Set the size of the swap to be 50% of RAM. Expected compression is 1:2. That
# means in an average case, 25% of RAM is consumed to provide the swap device.
# This results in the system being able to consume about 125% of RAM capacity to
# run processes. This value is inspired by Fedora which made zram the default
# swap mechanism. https://fedoraproject.org/wiki/Changes/SwapOnZRAM
[Service]
Environment=PERCENT=50
