# plinth/modules/ssh/data/usr/lib/systemd/system/ssh.service.d/freedombox.conf
[Unit]
# The service depends on users and groups defined in LDAP
After=nslcd.service
