# plinth/modules/syncthing/data/usr/lib/systemd/system/syncthing@syncthing.service.d/freedombox.conf
# The service depends on users and groups defined in LDAP
[Unit]
After=nslcd.service

# Revert to systemd defaults for the number of times the service can start in a
# given interval. Syncthing sets it to a much more strict value of 4 start in
# 60s which leads to functional tests failing when app is enabled/disabled,
# installed/uninstalled, and backed up/restored many times.
StartLimitIntervalSec=10s
StartLimitBurst=5
