# plinth/modules/upgrades/data/lib/systemd/system/freedombox-manual-upgrade.service
[Unit]
Description=Run unattended-upgrade once

[Service]
Type=oneshot
ExecStart=unattended-upgrade --verbose
KillMode=process
TimeoutStopSec=900
