# plinth/modules/upgrades/data/lib/systemd/system/freedombox-dist-upgrade.service
[Unit]
Description=Upgrade to new stable Debian release

[Service]
Type=oneshot
ExecStart=/usr/share/plinth/actions/upgrades dist-upgrade
KillMode=process
TimeoutStopSec=900
