# plinth/modules/coturn/data/usr/lib/systemd/system/coturn.service.d/freedombox.conf
[Service]
ExecStart=
ExecStart=/usr/bin/turnserver -c /etc/coturn/freedombox.conf --pidfile=/dev/null --log-file=-
ExecStartPost=
LockPersonality=yes
NoNewPrivileges=yes
PIDFile=
PrivateDevices=yes
PrivateMounts=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictRealtime=yes
SystemCallArchitectures=native
Type=simple
